# docs/agentes-reais/marvin/docs/com.marvin.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marvin.agent</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/agents/Desktop/claude-20x/agents/marvin/marvin_daemon.py</string>
        <string>start</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/agents/Desktop/claude-20x/agents/marvin</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/agents/Desktop/claude-20x/agents/marvin/logs/launchd_stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/agents/Desktop/claude-20x/agents/marvin/logs/launchd_stderr.log</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/agents/Desktop/claude-20x</string>
    </dict>
</dict>
</plist>